# PhoneClaw/supervisord.conf
[supervisord]
nodaemon=true
logfile=%(here)s/logs/supervisord.log
logfile_maxbytes=5MB
logfile_backups=2
pidfile=%(here)s/supervisord.pid
directory=%(here)s

[program:phoneclaw]
command=%(here)s/.venv/bin/python3 run.py
directory=%(here)s
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=%(here)s/logs/phoneclaw.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=%(here)s/logs/phoneclaw_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
environment=PYTHONUNBUFFERED="1"

[supervisorctl]
serverurl=unix:///data/data/com.termux/files/home/PhoneClaw/supervisor.sock

[unix_http_server]
file=%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
